# production/bitcoin.crontab
@reboot sleep 30 ; screen -dmS mainnet /litecoin/electrs/electrs-start-mainnet
@reboot sleep 60 ; /usr/local/bin/litecoind -testnet >/dev/null 2>&1
@reboot sleep 70 ; screen -dmS testnet /litecoin/electrs/electrs-start-testnet
@reboot sleep 80 ; /usr/local/bin/litecoind -signet >/dev/null 2>&1
@reboot sleep 90 ; screen -dmS signet /litecoin/electrs/electrs-start-signet
